# notify.conf
[program:tms]
command = /home/tms/tms_python/tms_env/bin/python /home/tms/tms_python/tms_env/bin/gunicorn "app:create_app()" --bind=0.0.0.0:8004 --access-logfile ../gunicorn-nf-access.log --error-logfile ../gunicorn-nf-error.log --capture-output --log-level debug --enable-stdio-inheritance
directory = /home/tms/notification/notification_system
user = tms

